# cmd/specula/systemd/specula.service
# Specula systemd unit — system-wide daemon, starts on boot.
# Canonical copy for `specula service install` (embedded). Also shipped under contrib/systemd/.

[Unit]
Description=Specula multi-protocol artifact proxy
Documentation=https://github.com/ivanzzeth/specula
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=specula
Group=specula
ExecStart=/usr/local/bin/specula --config /etc/specula/specula.yaml
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/specula /etc/specula
AmbientCapabilities=
CapabilityBoundingSet=

[Install]
WantedBy=multi-user.target
